# modbus-mqtt-bridge.service
[Service]
Type=simple
Restart=always
User=notroot
Group=notroot
WorkingDirectory=/home/notroot/src/mqtt-bridges
ExecStart=/usr/bin/env fish -ic 'npm run modbus'

[Install]
WantedBy=multi-user.target
